# system/madb-cache.timer
[Unit]
Description=Update metadata for MADB

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
